# restartnotify/files/unibas_restart_notifier_trigger.path
[Unit]
Description="Monitor the /boot directory for changes"

[Path]
PathModified=/boot
Unit=unibas_restart_notifier_trigger.service
TriggerLimitIntervalSec=1min

[Install]
WantedBy=multi-user.target